# config/examples/m12_37_intraday_auto_loop.systemd.timer.example
[Unit]
Description=Run Price Action Trader M12.37 readonly simulated dashboard every minute during user session

[Timer]
OnCalendar=*-*-* 09:30:00 America/New_York
OnUnitActiveSec=60s
Unit=m12_37_intraday_auto_loop.service

[Install]
WantedBy=timers.target

