# config/examples/m12_37_intraday_auto_loop.systemd.timer.example
[Unit]
Description=Start Price Action Trader M12.37 readonly simulated intraday session before US market open

[Timer]
OnCalendar=Mon..Fri *-*-* 09:25:00 America/New_York
Unit=m12_37_intraday_auto_loop.service
Persistent=true

[Install]
WantedBy=timers.target
